# chaturbate-dvr/rootfs/etc/services.d/chaturbate-dvr/run
#!/usr/bin/with-contenv bashio

# ==============================================================================
# Home Assistant Add-on: Chaturbate DVR
# ==============================================================================

# Get configuration values
USERNAME=$(bashio::config 'username')
FILENAME_PATTERN=$(bashio::config 'filename_pattern')
QUALITY=$(bashio::config 'quality')
USER_AGENT=$(bashio::config 'user_agent')
COOKIES=$(bashio::config 'cookies')
PROXY=$(bashio::config 'proxy')
LOG_LEVEL=$(bashio::config 'log_level')

# Determine output directory based on configuration
OUTPUT_LOCATION=$(bashio::config 'output_location')
if [ "$OUTPUT_LOCATION" = "custom" ]; then
    OUTPUT_DIR=$(bashio::config 'custom_output_path')
elif [ "$OUTPUT_LOCATION" = "media" ]; then
    OUTPUT_DIR="/media/chaturbate-dvr"
elif [ "$OUTPUT_LOCATION" = "config" ]; then
    OUTPUT_DIR="/config/chaturbate-dvr"
else
    # Default to share
    OUTPUT_DIR="/share/chaturbate-dvr"
fi

# Use ingress port (8080) for Home Assistant integration
PORT=8080

# Build command arguments
ARGS="-u ${USERNAME} -p ${PORT} -o ${OUTPUT_DIR}"

# Add filename pattern if configured
if bashio::config.has_value 'filename_pattern'; then
    ARGS="${ARGS} -f ${FILENAME_PATTERN}"
fi

# Add quality if configured
if bashio::config.has_value 'quality'; then
    ARGS="${ARGS} -q ${QUALITY}"
fi

# Add user agent if configured
if bashio::config.has_value 'user_agent'; then
    ARGS="${ARGS} -a ${USER_AGENT}"
fi

# Add cookies if configured
if bashio::config.has_value 'cookies'; then
    ARGS="${ARGS} -c ${COOKIES}"
fi

# Add proxy if configured
if bashio::config.has_value 'proxy'; then
    export HTTPS_PROXY="${PROXY}"
    export HTTP_PROXY="${PROXY}"
fi

# Add log level if configured
if bashio::config.has_value 'log_level'; then
    ARGS="${ARGS} -l ${LOG_LEVEL}"
fi

bashio::log.info "Starting Chaturbate DVR with arguments: ${ARGS}"
bashio::log.info "Web interface will be available through Home Assistant ingress"

# Start the application
exec chaturbate-dvr ${ARGS}
